# dump-database.timer
[Unit]
Description= restart every 10 minutes

[Timer]
OnCalendar=*:0/10:0
# OnCalendar=minutely

[Install]
WantedBy=timers.target
