# dump-database.service
[Unit]
Description=backup a database
After=network.target

[Service]
WorkingDirectory=/
User=root
ExecStart=/<path>/dump-database.sh
Type=oneshot
